# Logibooks.Backup/crontab
# Run backup daily at 2:00 AM
0 2 * * * /usr/local/bin/backup.sh >> /var/log/backup.log 2>&1

# Run backup every 6 hours (optional)
# 0 */6 * * * /usr/local/bin/backup.sh >> /var/log/backup.log 2>&1
